# dockerd
#!/bin/bash
#
# 
#
# chkconfig: 345 98 30
# description:  dockerd program.
# processname: dockerd

#start mysql
touch ~/$(date +%Y%m%d%M%S)

exit 0
